# testdata/fixtures/crontab/root.tab
MAILTO=ops@example.com
PATH=/usr/local/bin

*/5 * * * * kat render --theme Dracula
